# ansible/roles/docker_ingress_filter/templates/homelab-docker-ingress.service.j2
[Unit]
Description=Allowlist des ports Docker publiés du homelab
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-docker-ingress
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
