# services/osint/systemd/maigret-scan.service
[Unit]
Description=Maigret OSINT scan (profils famille)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/home/debian/osint/maigret/maigret-scan.sh
