# ansible/roles/scheduled_tasks/files/notify-fail@.service
[Unit]
Description=Alerte ntfy pour l'unité en échec %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelab-notify-fail %i
