# ansible/control-node/homelab-update-check.service
[Unit]
Description=Vérification des mises à jour disponibles sur le parc

[Service]
Type=oneshot
ExecStart=%h/.local/bin/homelab-update-check.sh
