# ansible/control-node/homelab-patch.timer
[Unit]
Description=Timer — fenêtre de patch hebdo (dimanche 04:30)

[Timer]
OnCalendar=Sun 04:30
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
